# deploy/digitalocean/systemd/faas-apid.service
[Unit]
Description=onebox-faas apid — control-plane API
After=network.target postgresql.service faas-cp.slice
Wants=faas-cp.slice
Requires=postgresql.service

[Service]
Type=simple
User=faas-apid
Group=faas
ExecStart=/opt/faas/bin/apid --config /etc/faas/apid.toml
Restart=on-failure
RestartSec=2s
Slice=faas-cp.slice
MemoryMax=256M
AmbientCapabilities=CAP_NET_BIND_SERVICE
EnvironmentFile=/etc/faas/sealed.env
# Session signing key (security review A4). LoadCredential surfaces
# /etc/faas/secrets/session.key inside the unit's credentials dir;
# Environment reads it back into FAAS_SESSION_KEY so apid's existing
# os.Getenv("FAAS_SESSION_KEY") loader path is unchanged. This keeps
# the key out of /etc/faas/sealed.env — every other control-plane
# daemon reads sealed.env but does NOT read this credential, so the
# 32-byte key never enters their environment.
LoadCredential=faas_session_key:/etc/faas/secrets/session.key
Environment=FAAS_SESSION_KEY=%d/faas_session_key
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadOnlyPaths=/etc/faas
ReadWritePaths=/var/log/faas /var/spool/faas

[Install]
WantedBy=multi-user.target
