# deploy/digitalocean/systemd/faas-cp.slice
[Unit]
Description=onebox-faas control plane (DO dev deployment)

[Slice]
MemoryMax=3G
