# deploy/ansible/roles/control_plane_service/files/faas-meterd.service
[Unit]
Description=onebox-faas meterd — metering, billing, and quota enforcement (spec §4.7)
Documentation=https://docs.DOMAIN/ops/meterd
After=network.target faas-cp.slice
Wants=faas-cp.slice

[Service]
Type=simple
User=faas-meterd
Group=faas
ExecStart=/usr/local/bin/meterd --config /etc/faas/meterd.toml
Restart=on-failure
RestartSec=2s

# meterd has no inbound gRPC; it dials schedd over a unix socket and
# PG over TCP. 192M covers the LISTEN goroutine + the rolling usage
# buffers for the hourly Stripe push (ADR-019).
Slice=faas-cp.slice
MemoryMax=192M

EnvironmentFile=/etc/faas/sealed.env

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadOnlyPaths=/etc/faas
ReadWritePaths=/var/log/faas

[Install]
WantedBy=multi-user.target
